# agent/systemd/seal-agent.service
[Unit]
Description=Seal host reconciliation
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/seal
ExecStart=/usr/local/bin/seal-agent /opt/seal
